# installer/l-panel.service
[Unit]
Description=L-Panel Web Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/l-panel
ExecStart=/opt/l-panel/venv/bin/python /opt/l-panel/run.py
Restart=always
RestartSec=5
User=root
Group=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
